# MQTT/mqtt-notifier.service
[Unit]
Description=MQTT notifier

[Service]
Type=oneshot
ExecStart=/usr/bin/mqttnotify.sh start
ExecStop=/usr/bin/mqttnotify.sh stop
StandardOutput=syslog

[Install]
WantedBy=multi-user.target

